# parking-scraper-report.service
[Unit]
Description=UNCC Parking Scraper Daily Report

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/uncc-parking-scraper
EnvironmentFile=/home/ubuntu/uncc-parking-scraper/.env
ExecStart=/home/ubuntu/uncc-parking-scraper/venv/bin/python -m src.daily_report
