# parking-scraper.timer
[Unit]
Description=Run parking scraper every 5 minutes on the dot

[Timer]
OnCalendar=*:0/5
AccuracySec=1s

[Install]
WantedBy=timers.target
