# parking-scraper.timer
[Unit]
Description=Run parking scraper every 15 minutes on the dot

[Timer]
OnCalendar=*:0/15
AccuracySec=1s

[Install]
WantedBy=timers.target
